# ProSim_Pi_UDP_Receiver_Stickshaker/my_server
### BEGIN INIT INFO
# Provides: my_server
# Required-Start: $remote_fs $syslog $network
# Required-Stop: $remote_fs $syslog $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Simple Web Server
# Description: Simple Web Server
### END INIT INFO

#! /bin/sh
# /etc/init.d/my_server

export HOME
echo "Starting My Server"
cd /home/pi/Documents/Flightsim/ProSim_Pi_UDP_Receiver_Stickshaker
sudo /usr/bin/python  Prosim_Pi_UDP_Receiver_Stickshaker_001.py 2>&1 &
exit 0

